# ansible/roles/mongodb/templates/mongod.service.j2
[Unit]
Description=High-performance, schema-free document-oriented database
After=network.target
Documentation=https://docs.mongodb.org/manual

[Service]
User=mongodb
Group=mongodb
ExecStart=/usr/bin/mongod --quiet --config {{mongodb_conf}}

[Install]
WantedBy=multi-user.target